# deploy/fake-email.service
[Unit]
Description=fake-email production stack (Docker Compose)
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/fake-email
ExecStartPre=/usr/bin/docker compose -f docker-compose.prod.yml pull --ignore-buildable
ExecStart=/usr/bin/docker compose -f docker-compose.prod.yml up -d --remove-orphans
ExecStop=/usr/bin/docker compose -f docker-compose.prod.yml down --timeout 30
ExecReload=/usr/bin/docker compose -f docker-compose.prod.yml up -d --remove-orphans

# Restart the compose stack if it exits unexpectedly
Restart=on-failure
RestartSec=10s

# Generous timeout for image builds / pulls
TimeoutStartSec=600
TimeoutStopSec=60

[Install]
WantedBy=multi-user.target
